# examples/cron.d/crontab3
# comment
* * * * * mblaschke echo "hello world!"
* * * * * mblaschke echo "hello" ; sleep 1 ; echo "world"
