# examples/cron.daily/hello-world3
#!/bin/sh

echo "Hello World! daily #1" >> /tmp/daily-3
