# examples/cron.daily/hello-world3
#!/usr/bin/env bash

echo "Hello World! daily #1" >> /tmp/daily-3
