# examples/cron.hourly/hello-world2
#!/usr/bin/env bash

echo "Hello World! hourly #1" >> /tmp/hourly-2
