# examples/cron.hourly/hello-world1
#!/usr/bin/env bash

echo "Hello World! hourly #1"
