# examples/crontab
# comment
* * * * * root sleep 5 && id >> /tmp/test-1
* * * * * root sleep 5 && id >> /tmp/test-2
@every 1m guest id >> /tmp/test-3
@every 1m guest env >> /tmp/test-4
